# sci-mathematics/sage-notebook-server/files/sage-notebook-server
#!/sbin/runscript

# Default values for variables defined in sage-notebook-server.conf
: ${user:=sage}
: ${accounts:=False}
: ${interface-localhost}
: ${port:=8000}
: ${secure:=False}
: ${require_login:=True}
: ${accounts:=False}
: ${timeout:=0}

user_home=$(eval cd ~${user} && pwd)
pidfile="${user_home}"/.sage/sage_notebook.sagenb/twistd.pid
command_args=( 
	-notebook
	accounts=${accounts}
	interface=${interface}
	open_viewer=False
	port=${port}
	require_login=${require_login}
	secure=${secure}
	timeout=${timeout}
)

: ${stdout:="${user_home}"/sagenb.log}
: ${stderr:="${user_home}"/sagenb.err}

depend() {
	need net
}

start() {
	# TODO: Add addtional code to detect different startup failures based on
	# standard output
	if [ ! -f "${user_home}"/.sage/sage_notebook.sagenb/users.pickle ] ; then
		eend 1 "Did not find user database - configure Sage's Notebook first:"

		echo
		einfo "- Configure the ${RC_SVCNAME}'s configuration file in: "
		einfo ""
		einfo "    /etc/conf.d/${RC_SVCNAME}.conf"
		einfo ""
		einfo "- Make sure you superuser and run the Notebook with the options"
		einfo "  you specified in the configuration file. Based on your current"
		einfo "  settings this would be:"
		einfo ""
		einfo "     su sage -c \"sage ${command_args[@]}\""
		einfo ""
		einfo "  Follow Sage's instructions and exit it after twistd is started"
		einfo "  by pressing Ctrl+C."
		einfo "- A multiuser server will require additional setup. Visit the"
		einfo "  following URL for details:"
		einfo ""
		einfo "    http://sagemath.org/doc/reference/sagenb/notebook/notebook_object.html"
		echo

		return 1
	fi

	ebegin "Starting Sage's Notebook server"
	start-stop-daemon \
		--start \
		--background \
		--stdout "${stdout}" \
		--stderr "${stderr}" \
		--pidfile "${pidfile}" \
		--user ${user} \
		--exec /usr/bin/sage \
		-- "${command_args[@]}"
	eend $? "Failed to start Sage's Notebook server"
}

stop() {
	ebegin "Stopping Sage's Notebook server"
	start-stop-daemon --stop --pidfile "${pidfile}" --user ${user}
	eend $? "Failed to stop Sage's Notebook server"
}
