# sci-mathematics/sage-notebook-server/files/sage-notebook-server
#!/sbin/runscript

extra_commands="reset"
description_reset="Sets/resets password for admin"

depend() {
	need net
}

start() {
	ebegin "Starting Sage's Notebook server"
	start-stop-daemon \
		--start \
		--background \
		--stdout ~sage/sagenb.log \
		--stderr ~sage/sagenb.err \
		--pidfile ~sage/sagenb.pid \
		--user sage \
		--exec /usr/bin/sage \
		-- -notebook \
		accounts=${accounts-False} \
		interface=${interface-localhost} \
		open_viewer=False \
		pidfile=/var/lib/sage/sagenb.pid \
		port=${port-8000} \
		require_login=${require_login-True} \
		secure=${secure-False} \
		timeout=${timeout-0}
	eend $? "Failed to start Sage's Notebook server"
}

stop() {
	ebegin "Stopping Sage's Notebook server"
	start-stop-daemon --stop --pidfile ~sagenb/sagenb.pid --user sage
	eend $? "Failed to stop Sage's Notebook server"
}

reset() {
	ebegin ""
	start-stop-daemon --start --pidfile ~sagenb/sagenb.pid --user sage --exec \
		/usr/bin/sage -- -notebook open_viewer=False reset=True
	eend $? ""
}
